# roles/glances/templates/supervisor_glances.conf.j2
[program:glances]
command=/usr/local/bin/glances -s

autostart=true
autorestart=true

user=root

#stdout_logfile=/var/logs/glances.log
#redirect_stderr = true